# app.cron
* * * * * echo "hello" > $HOME/cron.log 2>&1
